# scraper/crontab
# Chaos Scraper cron schedule
# Fires every hour 8:00–22:00, Mon–Sat.
# chaos_runner.sh sleeps a random 0–10 min before starting,
# so actual scraping begins at a random moment inside each hour.
# Logs: /var/log/scraper.log

# Format: minute hour dom month dow command
0 8-22 * * 1-6 /app/scripts/chaos_runner.sh >> /var/log/scraper.log 2>&1

# Empty line required at end of crontab file
