# scraper/crontab
# Scraper cron schedule
# Runs every hour during work hours (8:00-22:00), Monday-Saturday
# Logs are written to /var/log/scraper.log

# Format: minute hour day-of-month month day-of-week command
# Run at the start of every hour from 8 AM to 10 PM, Mon-Sat
0 8-22 * * 1-6 /app/scripts/scraper_runner.sh >> /var/log/scraper.log 2>&1

# Empty line required at end of crontab file
